# meta-balena-common/recipes-support/balena-config-vars/balena-config-vars/balena-engine-restarter.service
[Unit]
Description=balenaEngine restarter service
PartOf=config-json.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/restart-engine
