# meta-balena-common/recipes-containers/balena/balena/apps.slice
[Unit]
Description=balena user apps slice
Before=slices.target

[Slice]
CPUShares=16
