# meta-balena-common/recipes-containers/balena/balena/balena.slice
[Unit]
Description=balena critical services slice
Before=slices.target

[Slice]
CPUShares=262144
