# meta-balena-common/recipes-containers/balena/balena/scheduler-prios.service
[Unit]
Description=Set scheduling policies and priorities for critical balenaOS services

[Service]
ExecStart=/usr/lib/balena/set-os-scheduler-prios
